# data/case_123/sosreport-host0-2024-05-17-chmroof/usr/lib/systemd/system/multipathd.socket
[Unit]
Description=multipathd control socket
DefaultDependencies=no
ConditionPathExists=/etc/multipath.conf
ConditionKernelCommandLine=!nompath
ConditionKernelCommandLine=!multipath=off
ConditionVirtualization=!container
Before=sockets.target

[Socket]
ListenStream=@/org/kernel/linux/storage/multipathd

[Install]
WantedBy=sockets.target
