# data/case_123/sosreport-host0-2024-05-17-chmroof/usr/lib/systemd/system/realmd.service
[Unit]
Description=Realm and Dohost1 Configuration
Documentation=man:realm(8) man:realmd.conf(5)

[Service]
Type=dbus
BusName=org.freedesktop.realmd
ExecStart=/usr/libexec/realmd
