# data/case_123/sosreport-host0-2024-05-17-chmroof/usr/lib/systemd/system/rhcd.service
[Unit]
Description=Remote Host Configuration daemon
Documentation=https://github.com/redhatinsights/yggdrasil
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/rhcd

[Install]
WantedBy=multi-user.target
