# data/case_123/sosreport-host0-2024-05-17-chmroof/usr/lib/systemd/system/setroubleshootd.service
[Unit]
Description=SETroubleshoot daemon for processing new SELinux denial logs

[Service]
Type=dbus
BusName=org.fedoraproject.Setroubleshootd
ExecStart=/usr/sbin/setroubleshootd -f
User=setroubleshoot
LimitAS=1G
Nice=5
